# examples/dotfiles/systemd.host.user/frobnicate.service
[Service]
ExecStart=/usr/bin/frobnicate
