# examples/dotfiles/systemd.host.system/frobnicate@.service
[Service]
User=%i
ExecStart=/usr/bin/frobnicate
